# rootfs/etc/s6-overlay/s6-rc.d/moodle-init/up
/command/with-contenv /usr/local/bin/moodle-init.sh
