# rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
fdmove -c 2 1
php-fpm -F -O
